# systemd/spotify-swimmer.timer
# systemd/spotify-swimmer.timer
[Unit]
Description=Run Spotify Swimmer daily

[Timer]
# Base time at 5am, with up to 3 hours of random delay (runs between 5am-8am)
OnCalendar=*-*-* 05:00:00
RandomizedDelaySec=3h
Persistent=true

[Install]
WantedBy=timers.target
