# systemd/spotify-swimmer.service
# systemd/spotify-swimmer.service
[Unit]
Description=Spotify Swimmer - Download playlists for offline swimming
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=%h/.spotify-swimmer/bin/sync.sh
TimeoutStartSec=3600

[Install]
WantedBy=default.target
